# scripts/systemd/pcr-start.timer
[Unit]
Description=PCR Momentum – Daily Start Timer (Init + Stack)

[Timer]
OnCalendar=Mon-Fri *-*-* 08:00:00 Asia/Kolkata
Persistent=true
Unit=pcr-init.service

[Install]
WantedBy=timers.target
