# scripts/systemd/pcr-stack.target
[Unit]
Description=PCR Momentum Trading Bot – Cyclic Engine Stack
# Engines must NOT run unless init succeeded today.
Requires=pcr-init.service
After=pcr-init.service
Wants=pcr-scheduler.service pcr-data-pipeline.service pcr-background.service pcr-strategy.service pcr-order-exec.service pcr-health.service

[Install]
WantedBy=multi-user.target
