# scripts/systemd/pcr-stop.target
[Unit]
Description=PCR Momentum – Daily Stop Target (graceful shutdown of cyclic engines)

[Install]
WantedBy=multi-user.target
