# scripts/systemd/pcr-stack.target
[Unit]
Description=PCR Momentum Trading Bot – Cyclic Engine Stack
After=pcr-init.service
Requires=pcr-init.service

[Install]
WantedBy=multi-user.target
